# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#解决重启后磁盘变为只读导致不能保存配置的问题（每次开机重新挂载目录为rw）
mount -o remount rw /

#允许 Trojan 的不安全连接（不会检查远端主机所提供的 TLS 证书的有效性）
#其实仅需运行一次
sed -i 's/result.insecure = "0"/result.insecure = "1"/' /usr/share/shadowsocksr/subscribe.lua 

exit 0
